# scripts/crontabs/crontab.debian
NAME=jack@cwi.nl
MAILTO=jack@cwi.nl
PATH=/home/jack/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
50 11 * * * sh $HOME/src/ambulant/scripts/nightlyrun.debian.sh 
