# scripts/crontabs/crontab.linux
USER=jack@cwi.nl
MAILTO=jack@cwi.nl
30 00 * * * $HOME/src/ambulant/scripts/nightlybuild.linux.sh
