# scripts/crontabs/crontab.mac
NAME=jack@cwi.nl
MAILTO=jack@cwi.nl
PATH=/Users/jack/bin:/usr/local/bin:/usr/bin:/bin:/sbin
40 13 * * * $HOME/src/ambulant/scripts/nightlybuild.mac.sh
