# scripts/crontabs/crontab.mac
NAME=jack@cwi.nl
MAILTO=jack@cwi.nl
PATH=/Users/jack/bin:/usr/local/bin:/usr/bin:/bin:/sbin
20 0 * * * $HOME/src/ambulant/scripts/nightlybuild.mac.sh
10 14 * * * prlctl start WinXP-Nightly
