# scripts/crontabs/crontab.mac
NAME=jack@cwi.nl
MAILTO=jack@cwi.nl
PATH=/Users/jack/bin:/usr/local/bin:/usr/bin:/bin:/sbin
20 0 * * * sh $HOME/src/ambulant/scripts/nightlyrunmac.sh
20 2 * * * prlctl start WinXP-Nightly
