# scripts/crontabs/crontab.mac
30 00 * * * $HOME/src/ambulant/scripts/nightlybuild.mac.sh
