# scripts/crontabs/crontab.mac
NAME=jack@cwi.nl
MAILTO=jack@cwi.nl
50 11 * * * $HOME/src/ambulant/scripts/nightlybuild.mac.sh
