# scripts/crontabs/crontab.linux
MAILTO=ambulant-checkins@lists.sourceforge.net
30 00 * * * $HOME/src/ambulant/scripts/nightlybuild.linux.sh
