# scripts/nightlybuild/crontabs/crontab.debian
NAME=jack@cwi.nl
MAILTO=jack@cwi.nl
PATH=/home/jack/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
30 3 * * * sh $HOME/src/ambulant/scripts/nightlybuild/nightlyrun.linux.sh 
30 4 * * * sh $HOME/src/ambulant/scripts/nightlybuild/nightlyrun.debian.sh 
15 5 * * * ssh moes.sen.cwi.nl prlctl stop Ubuntu-1110-32bit
