# scripts/nightlybuild/crontabs/crontab.ubuntu-lts
NAME=jack@cwi.nl
MAILTO=jack@cwi.nl
PATH=/home/jack/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
30 2 * * * sh $HOME/src/ambulant/scripts/nightlybuild/nightlyrun.linux.sh 
30 3 * * * sh $HOME/src/ambulant/scripts/nightlybuild/nightlyrun.debian.sh 
15 4 * * * ssh nightlybuilds@maunakea.cwi.ta2 prlctl stop Ubuntu-1404-32bit
